# contrib/forked-daapd.conf
# forked-daapd

description	"forked-daapd"

start on runlevel [2345]
stop on runlevel [!2345]

expect fork
respawn

exec /usr/local/sbin/forked-daapd
